# mkosi.extra/usr/lib/systemd/system/lightdm.service.d/chown.conf
[Service]
ExecStartPre=/usr/bin/chown -R lightdm:lightdm /var/lib/lightdm
