# mkosi.extra/usr/lib/systemd/system/systemd-firstboot.service.d/etc.conf
[Service]
ExecStartPost=cp -rn /usr/share/factory/etc /
